# supervisord.stunnel.conf
[program:stunnel]
autorestart = true
command = /usr/bin/stunnel
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
