# supervisord.onvif-server.conf
[program:onvif-server]
autorestart = true
command = %(ENV_ONVIF_SERVER_HOME)s/onvifserver
directory = %(ENV_ONVIF_SERVER_HOME)s
environment = LD_LIBRARY_PATH="%(ENV_ONVIF_SERVER_HOME)s"
startretries = 0
startsecs = 10
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
